# infrastructure/systemd/preview-watch-uptime.service
[Unit]
Description=Preview Watch Uptime Probe
After=network.target

[Service]
Type=oneshot
ExecStart=/home/developerhiteshsaini/preview.watch/infrastructure/monitor/uptime_check.sh
User=root
Group=root
